# daemons/hyprpaper.service
[Unit]
Description=Wallpapers Daemon

[Service]
ExecStart=/usr/bin/hyprpaper
Restart=on-failure 

[Install]
WantedBy=default.target